# docker/crontab
0 */6 * * * . /etc/environment; /usr/local/bin/python /app/main.py >> /var/log/cron.log 2>&1
